# ansible/roles/baseline/templates/metrics.cron.j2
0 6 * * * root docker exec -u www-data {{ project_name }}_wp /usr/local/bin/send-wp-plugin-metrics
